# ansible/roles/autoupdate/templates/openfisca-web-api-autoupdate.timer.j2
[Unit]
Description=Run the Ansible playbook that installs the OpenFisca Web API.

[Timer]
OnCalendar={{ autoupdate_frequency }}
Persistent=true

[Install]
WantedBy=timers.target
