# fr.openfisca.org/legislation-explorer/legislation-explorer.service
[Unit]
Description=OpenFisca Legislation Explorer

[Service]
ExecStart=/usr/local/bin/npm start --prefix /home/openfisca/legislation-explorer
User=openfisca
Restart=always

[Install]
WantedBy=multi-user.target
