# ansible/roles/openfisca_api_fr/templates/cron.d/openfisca-api-fr.j2
# Run the Ansible playbook that installs the OpenFisca Web API.

0 */12 * * * root ansible-playbook --inventory ansible/inventories/local.yml ansible/site.yml
